# deployment/monitor.service
[Unit]
Description=Raspberry Monitor (port 10000)
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/monitor
Environment=MONITOR_PORT=10000
Environment=SPEEDTEST_INTERVAL_SEC=300
Environment=SPEEDTEST_TIMEOUT_SEC=60
Environment=TAILSCALE_CACHE_TTL_SEC=15
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Use pip-installed entry point
ExecStart=/usr/local/bin/raspberry-monitor
# Or use module: /usr/bin/python3 -m monitor
Restart=always
RestartSec=3

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
